# setup/service/pose-publisher.service
# Reference only — deploy.sh generates this file dynamically with the correct
# user and paths derived from the user@host argument.
#
# To install manually for user 'eikerd':
#   sudo cp this-file /etc/systemd/system/pose-publisher.service
#   sudo systemctl daemon-reload && sudo systemctl enable pose-publisher

[Unit]
Description=IMX500 Pose Publisher
After=network.target

[Service]
Type=simple
User=eikerd
WorkingDirectory=/home/eikerd/pose-publisher
ExecStart=/usr/bin/python3 /home/eikerd/pose-publisher/main.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
